# client/smolt.cron.monthly
# Runs the smolt checkin client

20 1 1 * * root /usr/bin/smoltSendProfile -c > /dev/null
